# com.frontier.linkedin-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.frontier.linkedin-worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/arjun/Documents/frontier_ventures/linkedin_mutual/venv/bin/python</string>
        <string>/Users/arjun/Documents/frontier_ventures/linkedin_mutual/worker.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/arjun/Documents/frontier_ventures/linkedin_mutual</string>

    <!-- Start automatically on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Wait 10s before restarting after a crash (avoid tight crash loops) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/arjun/Documents/frontier_ventures/linkedin_mutual/worker.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/arjun/Documents/frontier_ventures/linkedin_mutual/worker.error.log</string>

    <!-- Environment: inherit PATH so gcloud ADC and venv work correctly -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/arjun</string>
    </dict>

</dict>
</plist>
